# etc/skel/.config/autostart/kakaotalkset.desktop
[Desktop Entry]
Name=kakaotalkset
Comment=if not user... use this script
Exec=/usr/local/bin/kakaotalkset
Terminal=false
Type=Application
